# systemd/personal-ai-searcher.service
[Unit]
Description=personal-ai-searcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/personal-ai-searcher
EnvironmentFile=/opt/personal-ai-searcher/.env
ExecStart=/opt/personal-ai-searcher/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8020
Restart=always
RestartSec=5
User=deploy
Group=deploy
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
